# etc/udev/rules.d/91-usb-car-music.rules
ACTION=="add", SUBSYSTEM=="block", ENV{ID_FS_UUID}=="6630-278D", \
  RUN+="/home/context/bin/sync-music"
